# supervisord.d/condor-ce.conf
[program:condor-ce]
command=/usr/share/condor-ce/condor_ce_startup -f 
autorestart=true
